# deploy_configs/kapiqua_unicorn.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>kapiqua.unicorn</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>UserName</key>
    <string>proteus</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/EdwardData/Sites/kpiqa25/current/bin/unicorn</string>
        <string>-c</string>
        <string>/Users/EdwardData/Sites/kpiqa25/current/config/unicorn.rb</string>
        <string>-E</string>
        <string>production</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/EdwardData/Sites/kpiqa25/current</string>
  </dict>
</plist>